# debian/homevent.target
[Unit]
Description=The target to start all configured HomEvenT instances

[Install]
WantedBy=multi-user.target

