# debian/moat.target
[Unit]
Description=The target to start all configured MoaT instances

[Install]
WantedBy=multi-user.target

